# misc/supervisord.conf
[supervisord]
nodaemon=true

[program:mox]
command=mox
numprocs=1
autostart=true
autorestart=true
